# webserver.service
[Unit]
Description=musicazoo webserver service

[Service]
ExecStart=/usr/bin/python3 DIR/webserver.py

[Install]
WantedBy=default.target
