# wipe.service
[Unit]
Description=musicazoo wipe service

[Service]
ExecStart=DIR/wipe.sh
Restart=never

[Install]
WantedBy=default.target
